# services/syslog/run
#!/bin/execlineb -P
# Minimal syslog via s6-log

s6-log -d3 T /var/log/zigllm
