# services/syslog/run
#!/bin/execlineb -P
# s6 run script for syslog
s6-setuidgid root
/sbin/syslogd -n -O /var/log/messages
