# services/zigllm-ui/run
#!/bin/execlineb -P
# s6 run script for zigllm-ui
# Starts after zigllm-api is ready

s6-setuidgid root
/engine/zigllm-ui
