# services/arcis/finish
#!/bin/execlineb -P
# s6 finish script for arcis service
# Called when arcis exits. Logs exit code and restarts via s6 supervision.

import -u EXIT_CODE
foreground { s6-echo "arcis service exited with code ${EXIT_CODE} - s6 will restart" }
empty
