# services/mdev/run
#!/bin/execlineb -P
# s6 run script for mdev (device manager)
s6-setuidgid root
/sbin/mdev -s
